# deploy/supervisor/conf.d/inboxapp.conf
[program:inboxapp]
command=/srv/inboxapp/inboxapp-srv start
# stdout_logfile=/home/ubuntu/app_log.log
# redirect_stderr=true