# deploy/supervisor/conf.d/inboxapp.conf
[program:inboxapp]
command=/srv/inboxapp/inboxapp-srv debug
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log

