# deploy/supervisor/conf.d/mysqld.conf
[program:mysqld]
command=/usr/bin/mysqld_safe

