# deploy/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
stdout_logfile = /var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile = /var/log/supervisor/%(program_name)s-stderr.log
autorestart = true

